# ds5-edge-relay.service
[Unit]
Description=DualSense Edge HID Relay
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ds5-edge-relay
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=default.target
